# OpenBaton/scripts/HSS.service
[Unit]
Description=oai-epc-hss

[Service]
ExecStart=/root/openair-cn/SCRIPTS/run_hss

[Install]
WantedBy=multi-user.target
