# OpenBaton/scripts/MME.service
[Unit]
Description=oai-epc-mme

[Service]
ExecStart=/root/openair-cn/SCRIPTS/run_mme

[Install]
WantedBy=multi-user.target
