# OpenSourceMANO/scripts/SPGW.service
[Unit]
Description=oai-epc-spgw

[Service]
ExecStart=/root/openair-cn/SCRIPTS/run_spgw

[Install]
WantedBy=multi-user.target
